# chain_code_nchc/systemd/edit13-dispatcher.service
[Unit]
Description=Edit13_2800ITBLBM ITB-LBM dispatcher (32gpus@jp32 NCHC chain)
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/5.Re10595/Edit13_2800ITBLBM
ExecStart=/bin/bash %h/5.Re10595/Edit13_2800ITBLBM/chain_code_nchc/submit_dispatcher.sh
# crash / OOM / kill -> 自動重啟; 乾淨退出(STOP_CHAIN/STOP_DISPATCHER 收工, exit 0)-> 不復活
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
