# chain_code/systemd/edit11-dispatcher.service
[Unit]
Description=Edit11_Krank5600 GILBM dispatcher (net-best jp x partition auto-switch, never-idle)
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/s8313697/5.Re10595/Edit11_Krank5600
ExecStart=/bin/bash /home/s8313697/5.Re10595/Edit11_Krank5600/chain_code/submit_dispatcher.sh
# crash / OOM / kill -> 自動重啟; 乾淨退出(STOP_CHAIN/STOP_DISPATCHER 收工, exit 0)-> 不復活
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
