# chain_code/systemd/edit6-watchdog.timer
[Unit]
Description=Edit6_5600DNS health watchdog timer (every 10 min; semantic health beyond systemd Restart=on-failure)
[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
[Install]
WantedBy=timers.target
